# kismet-boot.service
# kismet-boot.service - Downloads JS file once at startup (one-shot)
[Unit]
Description=Download kismet boot script
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=toor
WorkingDirectory=/home/toor
ExecStart=/bin/bash -c 'curl -o /home/toor/kismet-boot.js https://raw.githubusercontent.com/c2-systems/rfd/refs/heads/main/kismet-boot.js && node kismet-boot.js'
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
